# distrib/build/sources/coretree/tree/etc/cron.daily/ssdtrim_cron
#!/bin/sh

if [-e /var/smoothwall/ssdtrim/supported]

	if [-e /var/smoothwall/ssdtrim/enabled]
		(
		echo "$(date -R)"
		fstrim -av
		echo ""
		) >> /var/log/ssdtrim.log
	
	else  
		exit
	fi

if [-e /var/smoothwall/ssdtrim/forced]
	(
		echo "$(date -R)"
		fstrim -av
		echo ""
		) >> /var/log/ssdtrim.log
	
fi
